# systemd/pi-fan-controller.service
[Unit]
Description=pi fan controller daemon

[Service]
Type=simple
ExecStart=/usr/bin/pi-fan-controller

[Install]
WantedBy=multiuser.agent