# upstart/pi-fan-controller.conf
# cpu fan controller

description "pi fan controller daemon"
author "Hordur Kristjansson (hoddikr@gmail.com)"

# Which GPIO pins to use in BCM numbers
env PIN1=2
env PIN2=3
env PIN3=4

# Set the CPU temperature thresholds
env TEMP1=55
env TEMP2=60
env TEMP3=65
env TEMP4=70
env TEMP5=75

# Set the interval in seconds to read the CPU temperature
env INTERVAL=5

# Default log level is LOG_NOTICE
# 0 Emergency
# 1 Alert
# 2 Critical
# 3 Error
# 4 Warning
# 5 Notice
# 6 Informational
# 7 Debug
env LOGLEVEL=6

# Set the time in seconds for the fan stopping grace period. This means that the controller will wait 
# for the defined amount of time until it completely stops the fan.
env STOPWAITTIME=120

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Essentially lets upstart know the process will detach itself to the background
# This option does not seem to be of great importance, so it does not need to be set.
expect fork

# Specify working directory
chdir /usr/bin

# Specify the process/command to start, e.g.
exec pi-fan-controller -i $INTERVAL -p1 $PIN1 -p2 $PIN2 -p3 $PIN3 -l $LOGLEVEL -st $STOPWAITTIME -t1 $TEMP1 -t2 $TEMP2 -t3 $TEMP3 -t4 $TEMP4 -t5 $TEMP5 -upstart